# etc/crontab
MAILTO=juangiordana@gmail.com

00 00 * * *         /srv/www/juangiordana.com.ar/bin/generate-sitemap.sh
